# os/systemd/karna-rtc-sync.service
[Unit]
Description=Sync hardware RTC from system clock after NTP synchronisation
Documentation=https://github.com/karna-project/KARNA-17
After=time-sync.target
Requires=time-sync.target

[Service]
Type=oneshot
ExecStart=/sbin/hwclock --systohc --utc
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=karna-rtc-sync

[Install]
WantedBy=time-sync.target
